# packer/nomad-server/nomad.mount
[Unit]
Description=Consul data drive

[Mount]
What=/dev/sda
Where=/mnt/nomad
Type=ext4
Options=defaults

[Install]
WantedBy=multi-user.target